# slack-copilot.service
[Unit]
Description=Slack Copilot Bot
After=network.target

[Service]
Type=simple
User=slack-copilot
Group=slack-copilot
WorkingDirectory=/opt/slack-copilot
Environment="PATH=/opt/slack-copilot/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/slack-copilot/.env
ExecStart=/usr/local/bin/uv run python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slack-copilot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/slack-copilot

[Install]
WantedBy=multi-user.target
